# supervisor/collectd.conf
[program:collectd]
user=root
directory=/
command=collectd -C /etc/collectd/collectd.conf -f
stdout_logfile=/var/log/supervisor/collectd.log
stderr_logfile=/var/log/supervisor/collectd.log
